# deploy/phoenixd/gamestr-arcade-wallet-backup.timer
[Unit]
Description=Daily encrypted backup of the Gamestr Arcade Phoenixd wallet

[Timer]
OnCalendar=*-*-* 03:17:00 UTC
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
